# gamestation-master/board/gamestation/fsoverlay/etc/udev/rules.d/99-sixpair.rules
#DRIVER=="usb", SUBSYSTEM=="usb", ATTR{idVendor}=="054c", ATTR{idProduct}=="0268", RUN+="/gamestation/scripts/udev-sixpair.sh"

